# overlays/configs/systemd/system/avahi-daemon.service.d/after-hostname.conf
[Unit]
After=set-hostname-and-banner.service
