# overlays/configs/systemd/user/plasma-kwin_wayland.service.d/override.conf
[Service]
Environment=KWIN_FORCE_SW_CURSOR=1
Environment=KWIN_COMPOSE=O2ES
