# overlays/configs/udev/rules.d/99-charger-status.rules
SUBSYSTEM=="power_supply", ACTION=="change", \
ENV{POWER_SUPPLY_STATUS}!="", \
RUN+="/usr/local/bin/battery-log.sh"
