# overlays/profile-tv-media-box/etc/systemd/system/kodi.service.d/handoff.conf
[Unit]
Conflicts=sddm.service
After=sddm.service
[Service]
ExecStop=-/usr/bin/kodi-send --action=Quit
TimeoutStopSec=5s