# overlays/profile-tv-media-box/etc/systemd/system/kodi.service.d/shutdown.conf
[Service]
# graceful stop (sddm/kodi never coexist on a profile, so the old handoff Conflicts/After
# were dead; this just asks kodi to quit cleanly on stop)
ExecStop=-/usr/bin/kodi-send --action=Quit
TimeoutStopSec=5s
