# overlays/configs/systemd/system/sddm.service.d/handoff.conf
[Unit]
Conflicts=kodi.service
After=kodi.service